# services/stack.target
[Unit]
Description=Start all service containers
